# docker/build/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=mysqld

[program:apache2]
command=apachectl -D FOREGROUND