# docker/build/conf/supervisord-prod.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=mysqld

[program:backend]
command=bash -c 'cd /root/backend && yarn start'

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true