# docker/build/conf/supervisord copy.conf
[supervisord]
nodaemon=true

[program:backend]
command=bash -c 'cd /root/backend && yarn start'

[program:mysqld]
command=mysqld

[program:apache2]
command=apachectl -D FOREGROUND